# scripts/systemd/manager-backup.timer
[Unit]
Description=Daily manager-komodo backup

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
