# scripts/systemd/manager-pghero-capture.service
[Unit]
Description=manager-komodo PGHero stats capture (history)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=ubuntu
Group=docker
WorkingDirectory=/opt/manager-komodo
ExecStart=/opt/manager-komodo/scripts/pghero-capture.sh
Environment=PATH=/usr/local/bin:/usr/bin:/bin
