# scripts/systemd/manager-pghero-refresh.timer
[Unit]
Description=Daily PGHero fleet config refresh

[Timer]
OnCalendar=*-*-* 03:40:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
